# ops/systemd/iic-sense-gdelt.service
[Unit]
Description=IIC sensing — GDELT adapter
After=network-online.target redis-server.service
Requires=redis-server.service

[Service]
Type=simple
User=ziwei-huang
WorkingDirectory=/home/ziwei-huang/TradingAgents/TradingAgents
EnvironmentFile=/home/ziwei-huang/TradingAgents/TradingAgents/.env
ExecStart=/home/ziwei-huang/miniconda3/bin/python -m tradingagents.sensing.adapters.gdelt
Restart=on-failure
RestartSec=30
MemoryMax=512M
CPUQuota=50%
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
