# ops/systemd/iic-triage.service
[Unit]
Description=IIC triage consumer — dedupe + salience + persistence
After=network-online.target redis-server.service
Requires=redis-server.service

[Service]
Type=simple
User=ziwei-huang
WorkingDirectory=/home/ziwei-huang/TradingAgents/TradingAgents
EnvironmentFile=/home/ziwei-huang/TradingAgents/TradingAgents/.env
ExecStart=/home/ziwei-huang/miniconda3/bin/python -m tradingagents.sensing.triage
Restart=on-failure
RestartSec=30
MemoryMax=2G
CPUQuota=100%
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
