# ops/systemd/iic-action-handler.service
[Unit]
Description=IIC-FORGE action handler — brief_actions sweep + dispatch
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ziwei-huang/TradingAgents
ExecStart=/home/ziwei-huang/TradingAgents/.venv/bin/forge action-handler run
Restart=on-failure
RestartSec=5
MemoryMax=1G
CPUQuota=50%
StandardOutput=journal
StandardError=journal
EnvironmentFile=/home/ziwei-huang/TradingAgents/.env

[Install]
WantedBy=multi-user.target
