# ops/systemd/iic-morning.service
[Unit]
Description=IIC-FORGE morning digest — oneshot
After=network.target

[Service]
Type=oneshot
User=ziwei-huang
WorkingDirectory=/home/ziwei-huang/TradingAgents/TradingAgents
ExecStart=/home/ziwei-huang/miniconda3/bin/python -m cli.main forge morning-digest now
StandardOutput=journal
StandardError=journal
EnvironmentFile=/home/ziwei-huang/TradingAgents/TradingAgents/.env
